# T1543.001 - Launch Agent.json
[
  {
    "Is the adversary maintaining persistence using Launch Agents? (TA0003 - Persistence)": {
      "Indicators": {
        "T1543.001 - Launch Agent": {
          "Creation or modification of new .plist files in user or system Launch Agent directories.": {
            "Data": "Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688 for process creations that involve writing files to known Launch Agent paths (e.g., paths resembling /Library/LaunchAgents). Use correlation analysis to link process creation events with subsequent file write events (if file auditing is available/mappable to Windows events) or immediate network connections (Zeek conn.log)."
          },
          "Launch Agents configured to execute files from suspicious directories (e.g., /tmp, /shared) or with suspicious names.": {
            "Data": "Windows Event ID 4688, Zeek conn.log",
            "Data Platform": "Endpoints, Servers, Network devices",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process executions where the command line or image path matches suspicious locations like /tmp/ or /shared/. Use frequency analysis on process names or paths to identify outliers. Correlate process executions with Zeek conn.log to detect immediate outbound connections."
          }
        }
      }
    }
  },
  {
    "Is the adversary attempting privilege escalation using Launch Agents? (TA0004 - Privilege-Escalation)": {
      "Indicators": {
        "T1543.001 - Launch Agent": {
          "Creation or modification of Launch Agents within a privileged user's profile or in system-wide directories with elevated permissions.": {
            "Data": "Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process creations under the context of privileged user accounts (e.g., Administrators, specific service accounts) or targeting system-wide configuration locations. Analyze the command line for execution patterns or parameters that are statistically unusual for that user."
          },
          "Execution of launchctl commands by processes or users not typically associated with system service management.": {
            "Data": "Windows Event ID 4688, Zeek conn.log",
            "Data Platform": "Endpoints, Servers, Network devices",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process creations where the command line contains launchctl or similar execution commands. Analyze the calling process's entropy or the execution frequency/timing using time series analysis to identify anomalies. Cross-reference with Zeek conn.log to detect network activity following execution."
          }
        }
      }
    }
  }
]