# T1543.001 - Launch Agent.json
[
  {
    "Is the adversary maintaining persistence using Launch Agents? (TA0003 - Persistence)": {
      "Indicators": {
        "T1543.001 - Launch Agent": {
          "Creation or modification of new .plist files in user or system Launch Agent directories.": {
            "Data": "Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688 for process creations that involve writing files to known Launch Agent paths (e.g., paths resembling /Library/LaunchAgents or ~/Library/LaunchAgents) using command-line utilities or scripting interpreters. Use correlation analysis to link process creation events with subsequent file write events (if file auditing is available/mappable to Windows events) or immediate network connections (Zeek conn.log)."
          },
          "Launch Agents configured to execute files from suspicious directories (e.g., /tmp, /shared) or with suspicious names.": {
            "Data": "Windows Event ID 4688; Zeek conn.log",
            "Data Platform": "Endpoints, Servers, Network devices",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process executions where the command line or image path matches suspicious locations like /tmp/ or /shared/. Use frequency analysis on process names or paths to identify outliers compared to a baseline. Correlate process executions with Zeek conn.log to detect immediate outbound connections."
          },
          "Command line arguments used by system utilities when creating or modifying files in Launch Agent directories that suggest setting persistence flags (e.g., RunAtLoad=true).": {
            "Data": "Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process executions (e.g., cmd.exe, powershell.exe, sh, bash, cp, mv, echo, tee) where the command line includes paths to Launch Agent directories (/Library/LaunchAgents, ~/Library/LaunchAgents, etc.) AND contains strings like 'RunAtLoad true' or 'KeepAlive true'. Use frequency analysis on command line arguments used in these specific contexts to identify rare or suspicious patterns not seen in typical system activity."
          },
          "Outbound network connections to unusual or known malicious destinations initiated by processes executed from Launch Agent paths.": {
            "Data": "Zeek conn.log; Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers, Network devices",
            "NAI": "Insert site-specific NAI here",
            "Action": "Identify process creations (Windows Event ID 4688) where the executable path is within standard user or system Launch Agent locations (/Library/LaunchAgents, ~/Library/LaunchAgents, etc.). Correlate these process start times with Zeek conn.log entries where the source IP matches the endpoint/server IP and the connection occurs shortly after the process started. Filter Zeek logs for connections to low-reputation or known malicious destination IPs using threat intelligence feeds. Use descriptive statistics on destination ports/protocols for these connections to identify unusual patterns compared to a baseline."
          }
        }
      },
      "version": "1.1",
      "last_updated": "2025-05-07"
    }
  },
  {
    "Is the adversary attempting privilege escalation using Launch Agents? (TA0004 - Privilege-Escalation)": {
      "Indicators": {
        "T1543.001 - Launch Agent": {
          "Creation or modification of Launch Agents within a privileged user's profile or in system-wide directories with elevated permissions.": {
            "Data": "Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process creations under the context of privileged user accounts (e.g., Administrators, specific service accounts) or targeting system-wide configuration locations (e.g., paths resembling /System/Library/LaunchAgents). Analyze the command line for execution patterns or parameters that are statistically unusual for that user or system context."
          },
          "Execution of launchctl commands by processes or users not typically associated with system service management.": {
            "Data": "Windows Event ID 4688; Zeek conn.log",
            "Data Platform": "Endpoints, Servers, Network devices",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process creations where the command line contains 'launchctl' or similar execution commands. Analyze the calling process's entropy or the execution frequency/timing using time series analysis to identify anomalies. Cross-reference with Zeek conn.log to detect network activity following execution, especially if originating from processes executing 'launchctl'."
          },
          "File creation or modification in system-wide Launch Agent directories (/System/Library/LaunchAgents) performed by processes running under standard user accounts.": {
            "Data": "Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process executions (e.g., cp, mv, text editors, scripting interpreters) that include paths like '/System/Library/LaunchAgents/' in their command line. For these events, check the 'SubjectUserName' or similar field to identify if the executing process is running under a standard (non-administrator, non-system) user account. Flag events where a standard user process attempts to interact with these system directories."
          },
          "Execution of launchctl commands (load, unload) specifically targeting .plist files within the system-wide Launch Agent directory (/System/Library/LaunchAgents) by processes running with standard user privileges.": {
            "Data": "Windows Event ID 4688",
            "Data Platform": "Endpoints, Servers",
            "NAI": "Insert site-specific NAI here",
            "Action": "Monitor Windows Event ID 4688. Filter for process executions where the command line contains 'launchctl load' or 'launchctl unload', and includes a path string matching '/System/Library/LaunchAgents/' and ending in '.plist'. For these events, check the 'SubjectUserName' to determine the user context. Flag events where a non-privileged user executes these specific commands against system Launch Agents. Use frequency analysis on the 'launchctl' command parameters and target files to identify unusual activity patterns for standard users."
          }
        }
      },
      "version": "1.1",
      "last_updated": "2025-05-07"
    }
  }
]