# config/asom.yml
# config/asom.yml
# Configuration for the Analytic Scheme of Maneuver (ASOM) generator.

# --- Input and Output Paths ---
# Directory containing the analytic plan JSON files.
analytic_plans_dir: "./techniques"

# Directory where the generated Excel reports will be saved.
output_dir: "./"

# --- Filename Settings ---
# Basename for the output files. A timestamp will be appended.
output_basename: "asom_test"

# Special filename for the 'Detect' tactic analytic plan.
detect_filename: "D3-D.json"

# --- Attack Chain Definition ---
# This section defines the structure of the ASOM.
# The script will process these tactics and techniques in the order they appear here.
# Tactic format: "TACTIC_ID - Tactic Name"
# Technique format: "TECHNIQUE_ID - Technique Name"
attack_chain:
  "D3-D - Detect":
    - "D3-NTA - Network Traffic Analysis"
    - "D3-PM - Platform Monitoring"
  "TA0001 - Initial Access":
    - "T1133 - External Remote Services"
  "TA0003 - Persistence":
    - "T1078 - Valid Accounts"
    - "T1053 - Scheduled Task/Job"
  "TA0008 - Lateral Movement":
    - "T1021 - Remote Services"
    - "T1570 - Lateral Tool Transfer"
  "TA0011 - Command and Control":
    - "T1071 - Application Layer Protocol"
  "TA0010 - Exfiltration":
    - "T1041 - Exfiltration Over C2 Channel"
    - "T1048 - Exfiltration Over Alternative Protocol"
    - "T1567 - Exfiltration Over Web Service"