# config/asom.yml
# config/asom.yml
# Configuration for the Analytic Scheme of Maneuver (ASOM) generator.

# --- Input and Output Paths ---
# Directories containing the analytic plan JSON files, categorized by framework.
# The script will search for .json files in each of these directories.
analytic_plan_dirs:
  enterprise: "techniques/attack/enterprise"
  ics: "techniques/attack/ics"
  mobile: "techniques/attack/mobile"
  d3fend: "techniques/d3fend"

# Directory where the generated Excel reports will be saved.
output_dir: "./"

# --- Filename Settings ---
# Basename for the output files. A timestamp will be appended.
output_basename: "asom_test"

# --- Detect Chain Definition ---
# Defines the "Detect" phase of the operation. These are processed first.
# Tactic format: "TACTIC_ID - Tactic Name"
# Technique format: "TECHNIQUE_ID - Technique Name"
detect_chain:
  "D3-D - Detect":
    - "D3-NTA - Network Traffic Analysis"
    - "D3-PM - Platform Monitoring"
  "D3-AI - Asset Inventory":
    - "D3-HCI - Hardware Component Inventory"
    - "D3-NNI - Network Node Inventory"
    - "D3-SWI - Software Inventory"
    - "D3-AVE - Asset Vulnerability Enumeration"
    - "D3-DI - Data Inventory"

# --- Attack Chain Definition ---
# Defines the MITRE ATT&CK phases of the operation. Processed after the detect chain.
attack_chain:
  "TA0001 - Initial Access":
    - "T1133 - External Remote Services"
  "TA0003 - Persistence":
    - "T1078 - Valid Accounts"
    - "T1053 - Scheduled Task/Job"
  "TA0008 - Lateral Movement":
    - "T1021 - Remote Services"
    - "T1570 - Lateral Tool Transfer"
  "TA0011 - Command and Control":
    - "T1071 - Application Layer Protocol"
  "TA0010 - Exfiltration":
    - "T1041 - Exfiltration Over C2 Channel"
    - "T1048 - Exfiltration Over Alternative Protocol"
    - "T1567 - Exfiltration Over Web Service"

# --- Output Columns ---
# Specify the exact columns to include in the final Excel reports.
# The order of columns in this list will be preserved in the output.
# If this list is empty or commented out, all columns will be included.
output_columns:
  - "CCIR Index"
  - "CCIR"
  # - "Tactic ID"
  # - "Tactic Name"
  - "Indicator Index"
  - "Indicator"
  # - "Technique ID"
  # - "Technique Name"
  - "Evidence Index"
  - "Evidence Description"
  - "Data Sources"
  # - "Data Platforms"
  - "NAI"
  - "Action"

# --- Formatting Settings ---
# Set approximate column widths in pixels for the Excel export.
# This is an approximation, as Excel's units are based on character width, not pixels.
# Columns not listed here will be auto-sized.
column_widths_pixels:
  "CCIR Index": 50
  "CCIR": 200
  "Tactic ID": 50
  "Tactic Name": 100
  "Indicator Index": 60
  "Indicator": 150
  "Technique ID": 75
  "Technique Name": 200
  "Evidence Index": 65
  "Evidence Description": 250
  "Data Sources": 200
  "Data Platforms": 150
  "NAI": 100
  "Action": 700